# deploy/nftably.service
[Unit]
Description=nftably - nftables web manager
Documentation=https://github.com/floreabogdan/nftably
After=network.target nftables.service
Wants=nftables.service

[Service]
Type=simple
User=nftably
Group=nftably

# Reading and applying the ruleset goes through netfilter (CAP_NET_ADMIN); the
# firewall log viewer reads the kernel ring buffer via dmesg, which needs
# CAP_SYSLOG when kernel.dmesg_restrict=1 (the default on most distros). nftably
# runs as its own unprivileged account with only these two capabilities — not
# full root — so a compromise of the web UI cannot, say, load a kernel module or
# change the clock.
AmbientCapabilities=CAP_NET_ADMIN CAP_SYSLOG
CapabilityBoundingSet=CAP_NET_ADMIN CAP_SYSLOG

ExecStart=/usr/bin/nftably server --db /var/lib/nftably/nftably.db
Restart=on-failure
RestartSec=5

# Hardening. nftably keeps its state under /var/lib/nftably and needs nothing
# else writable. ProtectSystem=strict makes the whole filesystem read-only for
# this unit, so the state dir has to be named explicitly.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/nftably
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes

[Install]
WantedBy=multi-user.target
